# pyLogger.conf
# pyLogger service 
# place in etc/init
description     "pyLogger"
author "github@thasler.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/henry/pyLogger/log.py
